# systemd/jitsi-calendar.service
[Unit]
Description=Launch Firefox kiosk for the active calendar meeting
After=graphical-session.target network-online.target
Wants=graphical-session.target network-online.target
ConditionUser=!root
StartLimitIntervalSec=5min
StartLimitBurst=3

[Service]
Type=oneshot
Environment=MOZ_ENABLE_WAYLAND=1
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=-/etc/JitsiRoomJoiner/.env
PassEnvironment=WAYLAND_DISPLAY DISPLAY XAUTHORITY XDG_RUNTIME_DIR DBUS_SESSION_BUS_ADDRESS
WorkingDirectory=/etc/JitsiRoomJoiner
ExecStartPre=/usr/bin/bash -c 'if [ -z "$WAYLAND_DISPLAY" ]; then echo "WAYLAND_DISPLAY not set; waiting for Wayland session" >&2; exit 1; fi'
ExecStartPre=/usr/bin/bash -c 'if [ -n "$WAYLAND_DISPLAY" ] && [ -n "$XDG_RUNTIME_DIR" ] && [ ! -S "$XDG_RUNTIME_DIR/$WAYLAND_DISPLAY" ]; then echo "Wayland socket $XDG_RUNTIME_DIR/$WAYLAND_DISPLAY not ready" >&2; exit 1; fi'
ExecStart=/usr/bin/env python3 /etc/JitsiRoomJoiner/jitsiRoomJoiner.py
TimeoutStartSec=90
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=graphical-session.target default.target
